# var/etc/cron.d/cron.disabled
* * * * * www-data /var/www/html/Siteworx/Cli/bin/app --cron >> /var/www/html/var/logs/cron.log

